# docker/rootfs/etc/services.d/rtpengine/run
#!/command/with-contenv sh
set -eu

RTPENGINE_BIN="$(command -v rtpengine || true)"
if [ -z "$RTPENGINE_BIN" ]; then
  RTPENGINE_BIN="/usr/sbin/rtpengine"
fi

PRIVATE_IP="${PRIVATE_IP:-}"
PUBLIC_IP="${PUBLIC_IP:-}"
MEDIA_IPV6="${MEDIA_IPV6:-0}"
MEDIA_FORCE_IPV4="${MEDIA_FORCE_IPV4:-}"

if [ -z "$MEDIA_FORCE_IPV4" ]; then
  if [ "$MEDIA_IPV6" = "1" ]; then
    MEDIA_FORCE_IPV4="0"
  else
    MEDIA_FORCE_IPV4="1"
  fi
fi

if [ -z "$PRIVATE_IP" ]; then
  PRIVATE_IP="$(ip -4 route get 1.1.1.1 2>/dev/null | awk '{for (i = 1; i <= NF; i++) if ($i == "src") {print $(i + 1); exit}}')"
fi

if [ -z "$PRIVATE_IP" ]; then
  PRIVATE_IP="127.0.0.1"
  echo "rtpengine: unable to detect PRIVATE_IP, falling back to ${PRIVATE_IP}" >&2
fi

if [ -z "$PUBLIC_IP" ]; then
  PUBLIC_IP="$PRIVATE_IP"
  echo "rtpengine: PUBLIC_IP not set, using ${PUBLIC_IP}" >&2
fi

if [ "$MEDIA_IPV6" = "1" ]; then
  echo "rtpengine: MEDIA_IPV6=1 (IPv6 media candidates enabled by policy), MEDIA_FORCE_IPV4=${MEDIA_FORCE_IPV4}" >&2
else
  echo "rtpengine: MEDIA_IPV6=0 (IPv4-only media candidates policy), MEDIA_FORCE_IPV4=${MEDIA_FORCE_IPV4}" >&2
fi

# 2223 is the NG UDP control socket used by Kamailio's rtpengine module.
# 7722 is the CLI TCP control socket used by rtpengine-ctl.
# Both sockets are bound to localhost and are not published externally.
exec "$RTPENGINE_BIN" \
  --foreground \
  --table=-1 \
  --interface="${PRIVATE_IP}!${PUBLIC_IP}" \
  --listen-ng=127.0.0.1:2223 \
  --listen-cli=127.0.0.1:7722 \
  --no-fallback \
  --port-min="${RTPENGINE_MIN_PORT}" \
  --port-max="${RTPENGINE_MAX_PORT}"
