# docker/rootfs/etc/services.d/healthcheck/run
#!/command/with-contenv sh
set -eu

exec python3 /opt/bitcall/healthcheck_server.py
